# utils/weimarnetz-basic-website/files/etc/uci-defaults/weimarnetz-basic-website-defaults
#!/bin/sh 

. /lib/functions.sh

uci_set uhttpd main index_page "index_basic.html"
uci_set uhttpd main error_page "/"
uci_commit
