# utils/weimarnetz-basic-website/files/etc/uci-defaults/weimarnetz-basic-website-defaults
#!/bin/sh 

. /lib/functions.sh

uci_set uhttpd main index_page "index2.html"
uci_commit
